# data/controls/IA-04.json
{
  "$schema": "../schema/control.schema.json",
  "id": "IA-04",
  "name": "Identifier Management",
  "family": "IA",
  "family_name": "Identification and Authentication",
  "control_class": "Technical",
  "description": "The organization manages user identifiers by: (i) uniquely identifying each user; (ii) verifying the identity of each user; (iii) receiving authorization to issue a user identifier from an appropriate organization official; (iv) issuing the user identifier to the intended party; (v) disabling the user identifier after [Assignment: organization-defined time period] of inactivity; and (vi) archiving user identifiers.",
  "supplemental_guidance": "Identifier management is not applicable to shared information system accounts (e.g., guest and anonymous accounts). FIPS 201 and Special Publications 800-73, 800- 76, and 800-78 specify a personal identity verification (PIV) credential for use in the unique identification and authentication of federal employees and contractors.",
  "enhancements": "(0) None.",
  "baseline_low": true,
  "baseline_moderate": true,
  "baseline_high": true,
  "joomla_id": 183,
  "nist_800_53": {
    "rev4": {
      "id": "IA-04",
      "name": "Identifier Management",
      "withdrawn": false,
      "incorporated_into": []
    },
    "rev5": {
      "id": "IA-04",
      "name": "Identifier Management",
      "description": "",
      "discussion": "",
      "related_controls": [],
      "baseline_low": true,
      "baseline_moderate": true,
      "baseline_high": true,
      "baseline_privacy": false,
      "new_in_rev5": false,
      "changes_from_rev4": "Removed control step to disable the identifier and associated parameter"
    }
  },
  "iso17799": [
    "11.2.3",
    "11.5.2"
  ],
  "cobit41": [
    "DS5.3",
    "DS5.4"
  ],
  "pci_dss_v2": [
    "8.5",
    "8.5.1",
    "8.5.5"
  ],
  "compliance_mappings": {
    "iso_27001_2022": [
      "A.5.16"
    ],
    "iso_27002_2022": [
      "5.16"
    ],
    "cobit_2019": [
      "DSS05"
    ],
    "pci_dss_v4": [
      "8.2"
    ],
    "nist_csf_2": [
      "PR.AA-01",
      "PR.AA-02"
    ],
    "cis_controls_v8": [
      "CIS 5",
      "CIS 5.5",
      "CIS 6.6"
    ],
    "soc2_tsc": [
      "CC6.1",
      "CC6.1-POF3",
      "CC6.1-POF4",
      "CC6.6",
      "CC6.6-POF2",
      "CC6.6-POF3"
    ],
    "finos_ccc": [
      "CCC-C11"
    ],
    "iso_42001_2023": [],
    "iec_62443": [],
    "asd_e8": [],
    "nis2": [],
    "apra_cps_234": [],
    "mas_trm": [
      "9"
    ],
    "uk_pra_fca": [],
    "bsi_grundschutz": [
      "ORP.4"
    ],
    "anssi": [
      "Hygiene.11",
      "Hygiene.32",
      "Hygiene.6",
      "Hygiene.7",
      "SecNumCloud.10.2"
    ],
    "osfi_b13": [
      "B-13.3.2"
    ],
    "finma_circular": [
      "IV.B.d(59)",
      "IV.B.d(60)"
    ],
    "gdpr": [
      "Art.32(1)(b)",
      "Art.5(1)(f)"
    ],
    "dora": [
      "Art.9(4)(c)",
      "Art.9(4)(d)"
    ],
    "bio2": [
      "5.16"
    ],
    "rbi_csf": [
      "Annex1.8",
      "ITGRCA.19"
    ],
    "fisc": [
      "FISC.T2"
    ],
    "lgpd_bcb": [],
    "hkma_tme1": [
      "TME1.8.1"
    ],
    "mlps_2": [
      "8.1.4.1"
    ],
    "dnb_good_practice": [
      "DNB.17.1",
      "DNB.17.2"
    ],
    "cra": [
      "CRA.I.2d"
    ],
    "swift_cscf": []
  },
  "attack_techniques": [
    {
      "id": "T1003",
      "name": "OS Credential Dumping",
      "tactics": [
        "credential-access"
      ],
      "mapping_type": "mitigates",
      "mapping_rationale": "Proper identifier management—including unique identification and timely disabling of inactive accounts—reduces the credential-dumping attack surface by ensuring that dormant or orphaned identifiers do not provide additional exploitable credential material."
    },
    {
      "id": "T1053",
      "name": "Scheduled Task/Job",
      "tactics": [
        "execution",
        "persistence",
        "privilege-escalation"
      ],
      "mapping_type": "mitigates",
      "mapping_rationale": "Identifier management with authorization controls ensures that scheduled tasks and jobs can only be created by properly identified and authorized users, preventing adversaries from scheduling persistent tasks under orphaned or shared identifiers."
    },
    {
      "id": "T1110",
      "name": "Brute Force",
      "tactics": [
        "credential-access"
      ],
      "mapping_type": "mitigates",
      "mapping_rationale": "Unique identifier management combined with timely deactivation of inactive accounts reduces the brute-force attack surface by eliminating dormant accounts that adversaries could target without triggering active-user alerting."
    },
    {
      "id": "T1213",
      "name": "Data from Information Repositories",
      "tactics": [
        "collection"
      ],
      "mapping_type": "mitigates",
      "mapping_rationale": "Proper identifier management ensures that access to information repositories is tied to verified, individual identifiers, enabling accurate access controls and audit trails that detect unauthorized data collection."
    },
    {
      "id": "T1528",
      "name": "Steal Application Access Token",
      "tactics": [
        "credential-access"
      ],
      "mapping_type": "mitigates",
      "mapping_rationale": "Identifier management ensures that application access tokens are bound to verified identities, enabling detection of token theft when tokens are used from systems or locations inconsistent with their assigned identifiers."
    },
    {
      "id": "T1530",
      "name": "Data from Cloud Storage",
      "tactics": [
        "collection"
      ],
      "mapping_type": "mitigates",
      "mapping_rationale": "Proper identifier management ensures cloud storage access is restricted to verified, authorized identifiers, preventing data collection through orphaned cloud accounts or shared credentials."
    },
    {
      "id": "T1537",
      "name": "Transfer Data to Cloud Account",
      "tactics": [
        "exfiltration"
      ],
      "mapping_type": "mitigates",
      "mapping_rationale": "Identifier management with timely deactivation prevents transfer of data to cloud accounts using orphaned or unauthorized identifiers, ensuring only verified users can establish cross-account data-transfer connections."
    },
    {
      "id": "T1543",
      "name": "Create or Modify System Process",
      "tactics": [
        "persistence",
        "privilege-escalation"
      ],
      "mapping_type": "mitigates",
      "mapping_rationale": "Identifier management ensures that system process creation and modification requires properly authorized identifiers, preventing adversaries from establishing persistent services using shared or orphaned account credentials."
    },
    {
      "id": "T1552",
      "name": "Unsecured Credentials",
      "tactics": [
        "credential-access"
      ],
      "mapping_type": "mitigates",
      "mapping_rationale": "Proper identifier lifecycle management reduces unsecured credential exposure by ensuring that disused identifiers and their associated credentials are archived or revoked, not left accessible for adversary exploitation."
    },
    {
      "id": "T1562",
      "name": "Impair Defenses",
      "tactics": [
        "defense-evasion"
      ],
      "mapping_type": "mitigates",
      "mapping_rationale": "Identifier management ensures that modifications to security controls require authenticated, authorized identifiers, preventing adversaries from impairing defenses using anonymous or shared accounts."
    },
    {
      "id": "T1563",
      "name": "Remote Service Session Hijacking",
      "tactics": [
        "lateral-movement"
      ],
      "mapping_type": "mitigates",
      "mapping_rationale": "Proper identifier management enables detection of remote session hijacking by tying sessions to unique identifiers, making it anomalous when a session is accessed from a different identity context."
    },
    {
      "id": "T1578",
      "name": "Modify Cloud Compute Infrastructure",
      "tactics": [
        "defense-evasion"
      ],
      "mapping_type": "mitigates",
      "mapping_rationale": "Identifier management ensures cloud compute modifications require verified identifiers, preventing adversaries from creating, modifying, or deleting cloud instances using orphaned or unauthorized cloud account credentials."
    },
    {
      "id": "T1602",
      "name": "Data from Configuration Repository",
      "tactics": [
        "collection"
      ],
      "mapping_type": "mitigates",
      "mapping_rationale": "Proper identifier management restricts access to network device configuration by requiring verified identifiers for management-protocol access, preventing unauthorized configuration extraction."
    },
    {
      "id": "T1003.005",
      "name": "Cached Domain Credentials",
      "tactics": [
        "credential-access"
      ],
      "mapping_type": "mitigates",
      "mapping_rationale": "Identifier management with timely account deactivation reduces cached domain credentials available for extraction by ensuring that dormant accounts do not accumulate stale cached logon material on workstations."
    },
    {
      "id": "T1003.006",
      "name": "DCSync",
      "tactics": [
        "credential-access"
      ],
      "mapping_type": "mitigates",
      "mapping_rationale": "Proper identifier management restricts directory replication rights to verified domain controller identifiers, making DCSync attacks detectable when replication requests originate from non-DC identifiers."
    },
    {
      "id": "T1021.001",
      "name": "Remote Desktop Protocol",
      "tactics": [
        "lateral-movement"
      ],
      "mapping_type": "mitigates",
      "mapping_rationale": "Identifier management ensures RDP access requires properly verified individual identifiers, preventing unauthorized remote desktop sessions through shared, default, or orphaned account credentials."
    },
    {
      "id": "T1021.005",
      "name": "VNC",
      "tactics": [
        "lateral-movement"
      ],
      "mapping_type": "mitigates",
      "mapping_rationale": "Proper identifier management restricts VNC access to verified identifiers, preventing adversaries from using shared or orphaned VNC credentials to establish unauthorized graphical remote control."
    },
    {
      "id": "T1053.002",
      "name": "At",
      "tactics": [
        "execution",
        "persistence",
        "privilege-escalation"
      ],
      "mapping_type": "mitigates",
      "mapping_rationale": "Identifier management ensures 'at' job creation requires authorized identifiers, preventing adversaries from scheduling persistent tasks through anonymous or orphaned batch-scheduling accounts."
    },
    {
      "id": "T1053.005",
      "name": "Scheduled Task",
      "tactics": [
        "execution",
        "persistence",
        "privilege-escalation"
      ],
      "mapping_type": "mitigates",
      "mapping_rationale": "Proper identifier management ensures scheduled task creation requires verified identifiers, enabling detection when tasks are created under orphaned or unauthorized account credentials."
    },
    {
      "id": "T1098.007",
      "name": "Additional Local or Domain Groups",
      "tactics": [
        "persistence",
        "privilege-escalation"
      ],
      "mapping_type": "mitigates",
      "mapping_rationale": "Identifier management prevents unauthorized group membership changes by ensuring group modifications require verified administrator identifiers, detecting adversary attempts to add themselves to privileged groups."
    },
    {
      "id": "T1110.001",
      "name": "Password Guessing",
      "tactics": [
        "credential-access"
      ],
      "mapping_type": "mitigates",
      "mapping_rationale": "Unique identifier management with inactive-account disabling reduces password-guessing targets by eliminating dormant accounts from the authentication-attempt surface."
    },
    {
      "id": "T1110.002",
      "name": "Password Cracking",
      "tactics": [
        "credential-access"
      ],
      "mapping_type": "mitigates",
      "mapping_rationale": "Proper identifier lifecycle management reduces password-cracking value by ensuring that credentials associated with disabled identifiers cannot be used for authentication even if hashes are cracked offline."
    },
    {
      "id": "T1110.003",
      "name": "Password Spraying",
      "tactics": [
        "credential-access"
      ],
      "mapping_type": "mitigates",
      "mapping_rationale": "Identifier management with unique-user identification makes password-spraying detectable by enabling per-account failed-authentication tracking that identifies systematic credential-testing patterns."
    },
    {
      "id": "T1110.004",
      "name": "Credential Stuffing",
      "tactics": [
        "credential-access"
      ],
      "mapping_type": "mitigates",
      "mapping_rationale": "Proper identifier management reduces credential-stuffing effectiveness by ensuring that identifiers are unique to the organization, limiting the applicability of leaked credentials from external breaches."
    },
    {
      "id": "T1213.001",
      "name": "Confluence",
      "tactics": [
        "collection"
      ],
      "mapping_type": "mitigates",
      "mapping_rationale": "Identifier management ensures Confluence access is tied to verified identifiers, enabling detection of unauthorized data collection through access-pattern analysis against individual account baselines."
    },
    {
      "id": "T1213.002",
      "name": "Sharepoint",
      "tactics": [
        "collection"
      ],
      "mapping_type": "mitigates",
      "mapping_rationale": "Proper identifier management ensures SharePoint access requires verified identifiers, enabling access auditing and detection of anomalous document-library access patterns."
    },
    {
      "id": "T1213.004",
      "name": "Customer Relationship Management Software",
      "tactics": [
        "collection"
      ],
      "mapping_type": "mitigates",
      "mapping_rationale": "Identifier management ensures CRM access requires verified identifiers, preventing unauthorized access to customer data through shared, default, or orphaned account credentials."
    },
    {
      "id": "T1213.005",
      "name": "Messaging Applications",
      "tactics": [
        "collection"
      ],
      "mapping_type": "mitigates",
      "mapping_rationale": "Proper identifier management ensures messaging application access requires verified identifiers, enabling detection of unauthorized message access through per-user activity baselines."
    },
    {
      "id": "T1547.006",
      "name": "Kernel Modules and Extensions",
      "tactics": [
        "persistence",
        "privilege-escalation"
      ],
      "mapping_type": "mitigates",
      "mapping_rationale": "Identifier management ensures that kernel module installation requires verified administrative identifiers, preventing adversaries from loading rootkits or malicious extensions under unauthorized credentials."
    },
    {
      "id": "T1550.001",
      "name": "Application Access Token",
      "tactics": [
        "defense-evasion",
        "lateral-movement"
      ],
      "mapping_type": "mitigates",
      "mapping_rationale": "Proper identifier management supports application-access-token validation by ensuring tokens are bound to verified identifiers, enabling detection when tokens are used outside their assigned identity context."
    },
    {
      "id": "T1552.005",
      "name": "Cloud Instance Metadata API",
      "tactics": [
        "credential-access"
      ],
      "mapping_type": "mitigates",
      "mapping_rationale": "Identifier management with proper cloud identity controls restricts metadata API credential access to verified instance identifiers, preventing unauthorized credential harvesting from cloud environments."
    },
    {
      "id": "T1578.001",
      "name": "Create Snapshot",
      "tactics": [
        "defense-evasion"
      ],
      "mapping_type": "mitigates",
      "mapping_rationale": "Proper identifier management ensures cloud snapshot creation requires verified identifiers, preventing adversaries from exfiltrating data by creating unauthorized snapshots of production volumes."
    },
    {
      "id": "T1578.002",
      "name": "Create Cloud Instance",
      "tactics": [
        "defense-evasion"
      ],
      "mapping_type": "mitigates",
      "mapping_rationale": "Identifier management ensures cloud instance creation requires verified identifiers, preventing adversaries from spinning up unauthorized compute resources for cryptomining or pivot operations."
    },
    {
      "id": "T1578.003",
      "name": "Delete Cloud Instance",
      "tactics": [
        "defense-evasion"
      ],
      "mapping_type": "mitigates",
      "mapping_rationale": "Proper identifier management ensures cloud instance deletion requires verified identifiers, preventing adversaries from destroying organizational resources through unauthorized teardown operations."
    },
    {
      "id": "T1602.001",
      "name": "SNMP (MIB Dump)",
      "tactics": [
        "collection"
      ],
      "mapping_type": "mitigates",
      "mapping_rationale": "Identifier management ensures SNMP access requires authenticated identifiers (SNMPv3), preventing unauthorized MIB queries through community-string-based access that lacks individual accountability."
    },
    {
      "id": "T1602.002",
      "name": "Network Device Configuration Dump",
      "tactics": [
        "collection"
      ],
      "mapping_type": "mitigates",
      "mapping_rationale": "Proper identifier management ensures network device management requires verified individual identifiers, preventing unauthorized configuration extraction through shared or default management credentials."
    }
  ],
  "metadata": {
    "last_reviewed": "2026-02-13",
    "review_notes": "",
    "mapping_status": "complete"
  }
}
